# deferred_wifi
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: MIT

START=99

boot() {
        wifi
}

start() {
        wifi
}

restart() {
        wifi
}

stop() {
        return 0
}

